# templates/postgresql.service.j2
.include /lib/systemd/system/{{ postgresql_service_name }}.service

[Service]
Environment=PGPORT={{ postgresql_connection.port }}
Environment=PGDATA={{ postgresql.data_directory }}
